# docker/etc/cron.d/myddleware
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
SHELL=/bin/bash

# Run Myddleware cron every minute - output directly to log (rsyslog not available in Docker)
*/1 * * * * www-data /usr/local/bin/myddleware-cron.sh >> /var/log/myddleware-cron.log 2>&1

# Monitor running jobs every minute - disabled until myddleware-monitor.sh is implemented
# * * * * * root /usr/local/bin/myddleware-monitor.sh >> /var/log/myddleware-monitor.log 2>&1
