# docker/etc/cron.d/myddleware
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
SHELL=/bin/bash

*/5 * * * * www-data /usr/local/bin/myddleware-cron.sh >> /var/log/myddleware-cron.log 2>&1
* * * * * root /usr/local/bin/myddleware-monitor.sh >> /var/log/myddleware-monitor.log 2>&1
