# docker/etc/crontab
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
*/5 * * * * www-data myddleware-cron.sh 2>&1
* * * * * root bash myddleware-monitor.sh 2>&1
